# systemd/.config/systemd/user/opencode-serve.service
[Unit]
Description=OpenCode AI Assistant Server
After=network.target

[Service]
Type=simple
# Run opencode through mise to ensure Node.js and other tools are available
ExecStart=/home/linuxbrew/.linuxbrew/bin/mise exec -- /home/linuxbrew/.linuxbrew/bin/opencode serve --port 1337
Restart=on-failure
RestartSec=5
Environment=HOME=%h
Environment=USER=%u
# mise needs these to find its config
Environment=MISE_DATA_DIR=%h/.local/share/mise
Environment=MISE_CONFIG_DIR=%h/.config/mise
# CRITICAL: Add mise shims to PATH so opencode's subprocesses (node, etc) route through mise
# This fixes project tools not being found - shims auto-detect project mise.toml via PWD
Environment=PATH=%h/.local/share/mise/shims:/home/linuxbrew/.linuxbrew/shims:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin

# Ensure proper runtime directory
RuntimeDirectory=%t/opencode
RuntimeDirectoryMode=0755

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode

[Install]
WantedBy=default.target
